# contrib/linux-packaging/archlinux/neucoind.service
[Unit]
Description=Neucoin daemon
After=network.target

[Service]
ExecStart=neucoind -logtimestamps

[Install]
WantedBy=default.target
